# systemd/habit.service
[Unit]
Description=HabitFlow Habit Service
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/opt/habitflow/habit-service
EnvironmentFile=/opt/habitflow/habit-service/.env
ExecStart=/opt/habitflow/.venv/bin/uvicorn app:app --host 0.0.0.0 --port 8002
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
